# patch/files/etc/uci-defaults/99-network
uci set network.lan.ipaddr='10.11.12.13'
uci commit network
