# patch/files/etc/uci-defaults/99-timezone
uci set system.@system[0].timezone='WIB-7'
uci set system.@system[0].zonename='Asia/Jakarta'
uci set system.@system[0].conloglevel='5'
uci set system.@system[0].cronloglevel='9'
uci delete system.ntp.server
uci add_list system.ntp.server='0.id.pool.ntp.org'
uci add_list system.ntp.server='1.id.pool.ntp.org'
uci add_list system.ntp.server='2.id.pool.ntp.org'
uci commit system
/etc/init.d/system restart
